# lede6.6/target/linux/amlogic/base-files/etc/uci-defaults/dhcp-proto-lan
#!/bin/sh
#
# Set dhcp proto on lan

uci -q batch <<-EOF >/dev/null
	delete network.lan.ipaddr
	delete network.lan.netmask
	set network.lan.proto='dhcp'
	commit network
EOF

rm -f /tmp/luci-indexcache
exit 0
